# sysfiles/desktop/usr/lib/systemd/system/suspendfix.service
[Unit]
Description=Fix suspend with Gigabyte board
[Service]
Type=oneshot
ExecStart=/usr/bin/suspendfix.sh
[Install]
WantedBy=multi-user.target
