# rootfs/usr/lib/systemd/system/scx-lavd.service
[Unit]
Description=Enables LAVD scheduler on boot

[Service]
Type=oneshot
ExecStart=/usr/bin/scx_lavd

[Install]
WantedBy=multi-user.target
